# launchagents/com.openclaw.sequence.simple-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
<key>Label</key>
<string>com.openclaw.sequence.simple-worker</string>
<key>ProgramArguments</key>
<array>
<string>/bin/zsh</string>
<string>-lc</string>
<string>export PATH="/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:$PATH"; set -a; source $HOME/.assistantplatform/.env; set +a; cd $HOME/.openclaw/projectdir/sequence-service; source .venv/bin/activate; exec python3 simple_worker.py</string>
</array>
<key>RunAtLoad</key>
<true/>
<key>KeepAlive</key>
<true/>
<key>StandardOutPath</key>
<string>/tmp/sequence_simple_worker.out.log</string>
<key>StandardErrorPath</key>
<string>/tmp/sequence_simple_worker.err.log</string>
</dict>
</plist>
